# deploy/supervisor.ini
[program:dj_renesola]
user=marc
command=/home/marc/envs/renesola/bin/python /home/marc/python_test_apps/renesola/manage.py runserver 0.0.0.0:8102
directory=/home/marc/python_test_apps/renesola/
autostart=true
autorestart=true
logfile=log.txt


stdout_logfile = /home/marc/python_test_apps/renesola/server.log

redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
