# renesola/kombu_workers/deploy/sfgeek.net.supervisor.ini
[program:dj_renesola_kombu]
user=marc
command=/home/marc/envs/renesola/bin/python /home/marc/python_test_apps/renesola/renesola/kombu_workers/worker.py
directory=/home/marc/python_test_apps/renesola/renesola/kombu_workers/
autostart=true
autorestart=true

stdout_logfile = /home/marc/python_test_apps/renesola/renesola/kombu_workers/worker.log

redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
