# renesola/kombu_workers/deploy/supervisor.ini
[program:dj_renesola_kombu]
user=ubuntu
command=python /var/www/crm/renesola/kombu_workers/worker.py
directory=/var/www/crm/renesola/kombu_workers/
autostart=true
autorestart=true

stdout_logfile = /var/www/crm/renesola/kombu_workers/worker.log

redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
